# com.mcpcompanion.services.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcpcompanion.services</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mark/Desktop/MCP STUFF/mcp-macos-companion/start_services.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/mark/.mcp/logs/launchd.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mark/.mcp/logs/launchd.err</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/mark/Desktop/MCP STUFF/mcp-macos-companion</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>